# sci-geosciences/kplex/files/kplex.initd
#!/sbin/openrc-run
# Copyright 2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="kplex daemon"
description="An NMEA-0183 Multiplexer"
pidfile="/run/${RC_SVCNAME}.pid"
command=/usr/bin/kplex
command_args="-p ${pidfile} ${kplex_args}"
command_args_background="-o mode=background"

depend() {
	need net
	use logger
}
